# manis/Daemon/Launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.manis.Daemon</string>
    <key>MachServices</key>
    <dict>
        <key>com.manis.Daemon</key>
        <true/>
    </dict>
    <key>Program</key>
    <string>/Library/Daemons/com.manis.Daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Daemons/com.manis.Daemon</string>
    </array>
    <key>AssociatedBundleIdentifiers</key>
    <string>com.manis.app</string>
</dict>
</plist>